# config/crontab
0 22 * * * DISPLAY=:0.0 /usr/bin/redshift -O 4000
